# redhawk-container/docker-redhawk/Dockerfiles/files/supervisord-device.conf
[program:nodeBooter]
autostart=true
redirect_stderr=true
command=/bin/bash -lc "wait-for-domain nodeBooter -d /nodes/$NODENAME/DeviceManager.dcd.xml"
startsecs=35
priority=10

[eventlistener:nodeBooter_exit]
command=/usr/bin/kill_supervisor.py
process_name=nodeBooter
events=PROCESS_STATE_EXITED
