# redhawk/src/etc/profile.d/redhawk.sh.in
#
# This file is protected by Copyright. Please refer to the COPYRIGHT file 
# distributed with this source distribution.
# 
# This file is part of REDHAWK core.
# 
# REDHAWK core is free software: you can redistribute it and/or modify it under 
# the terms of the GNU Lesser General Public License as published by the Free 
# Software Foundation, either version 3 of the License, or (at your option) any 
# later version.
# 
# REDHAWK core is distributed in the hope that it will be useful, but WITHOUT 
# ANY WARRANTY; without even the implied warranty of MERCHANTABILITY or FITNESS 
# FOR A PARTICULAR PURPOSE.  See the GNU Lesser General Public License for more
# details.
# 
# You should have received a copy of the GNU Lesser General Public License 
# along with this program.  If not, see http://www.gnu.org/licenses/.
#

if [ "X$OSSIEHOME" = X ] ; then
    export OSSIEHOME=@prefix@
    if [ "X$PATH" != X ]; then
        PATH=:${PATH}
    fi
    export PATH=@bindir@${PATH}
    if [ "X$LD_LIBRARY_PATH" != X ]; then
        LD_LIBRARY_PATH=:${LD_LIBRARY_PATH}
    fi
    export LD_LIBRARY_PATH=@libdir@${LD_LIBRARY_PATH}
    if [ "X$PYTHONPATH" != X ]; then
        PYTHONPATH=:${PYTHONPATH}
    fi
    export PYTHONPATH=${OSSIEHOME}/lib64/python:${OSSIEHOME}/lib/python${PYTHONPATH}
fi
