# redhawk/src/etc/profile.d/redhawk-sdrroot.csh.in
if (! $?SDRROOT ) then
    setenv SDRROOT @sdrroot@
endif

