# redhawk/src/etc/profile.d/redhawk.csh.in
if (! $?OSSIEHOME) then
    setenv OSSIEHOME @prefix@

    if ($?PATH) then
        setenv PATH @bindir@:${PATH}
    else
        setenv PATH @bindir@
    endif

    if ($?PYTHONPATH) then
        setenv PYTHONPATH $OSSIEHOME/lib64/python:$OSSIEHOME/lib/python:${PYTHONPATH}
    else
        setenv PYTHONPATH $OSSIEHOME/lib64/python:$OSSIEHOME/lib/python
    endif
endif
if (! $?KUBECONFIG) then
    setenv KUBECONFIG $OSSIEHOME/.kube/config
endif
if (! $?AWS_CONFIG_FILE) then
    setenv AWS_CONFIG_FILE $OSSIEHOME/aws/config
endif
if (! $?AWS_SHARED_CREDENTIALS_FILE) then
    setenv AWS_SHARED_CREDENTIALS_FILE $OSSIEHOME/aws/credentials
endif
if (! $?AWS_PROFILE) then
    setenv AWS_PROFILE redhawk
endif
if (! $?LD_LIBRARY_PATH) then
    setenv LD_LIBRARY_PATH @libdir@
else
    setenv LD_LIBRARY_PATH @libdir@:${LD_LIBRARY_PATH}
endif
